# etc/systemd/system/docker.service.d/override.conf
[Service]

# Environment="HTTP_PROXY=http://proxy.example.com:80/"
# Environment="HTTPS_PROXY=http://proxy.example.com:80/"
# Environment="NO_PROXY=localhost,127.0.0.1,private.docker.registry.com"

Environment="all_proxy=socks5://127.0.0.1:1080"

