# etc/systemd/system/docker.service.d/override.conf
[Service]

## uncomment for low-speed network or private registry:
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --max-concurrent-downloads 1 --containerd=/run/containerd/containerd.sock
#ExecStart=/usr/bin/dockerd -H fd:// -H tcp://0.0.0.0:2375 --insecure-registry=private.docker.registry.com --insecure-registry=private.kube.registry.com --registry-mirror=http://private.kube.registry.com:6000

## uncomment for blocked/private networks
Environment="ALL_PROXY=socks5h://127.0.0.1:1080"
Environment="NO_PROXY=localhost,127.0.0.1"
# Environment="NO_PROXY=localhost,127.0.0.1,private.docker.registry.com"
# Environment="HTTP_PROXY=http://proxy.example.com:80/"
# Environment="HTTPS_PROXY=http://proxy.example.com:80/"
